# infra/ddns/cloudflare-ddns.service
[Unit]
Description=Cloudflare DDNS updater for local.kydyrov.dev
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/cloudflare-ddns
ExecStart=/usr/bin/python3 /usr/local/lib/devops-platform-ddns/cloudflare_ddns.py
